# .config/autostart/quake.desktop
[Desktop Entry]
Exec=/usr/bin/guake
Name=guake
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
